# dot_config/systemd/user/secret-monitor.service
[Unit]
Description=Secret access monitor (enriches KeePassXC prompts)
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.config/sway/bin/secret-monitor.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
