# dot_config/systemd/user/swaync.service
[Unit]
Description=Sway notification center
After=graphical-session.target
PartOf=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
Type=simple
ExecStart=/usr/bin/swaync
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=30
StartLimitBurst=3

[Install]
WantedBy=graphical-session.target
