# dot_config/systemd/user/systembus-notify.service
[Unit]
Description=Forward system D-Bus notifications to desktop
After=graphical-session.target

[Service]
ExecStart=/usr/bin/systembus-notify
Restart=on-failure

[Install]
WantedBy=graphical-session.target
